# dev/hostapd.service
[Service]
RestartSec=5
Restart=on-failure
